# dist/systemd/system/crau-event-broker.socket
[Unit]
Description=Event broker socket for crypto-auditing project

[Socket]
ListenStream=/run/crypto-auditing/audit.sock
# SocketUser=crypto-auditing
# SocketGroup=crypto-auditing
# SocketMode=0660

[Install]
WantedBy=sockets.target
