# dist/systemd/system/crypto-auditing-agent.service
[Unit]
Description=Event collector agent for crypto-auditing project
StartLimitInterval=10s
StartLimitBurst=5
After=local-fs.target systemd-tmpfiles-setup.service
Before=sysinit.target shutdown.target

[Service]
ExecStart=/usr/bin/crypto-auditing-agent
KillSignal=SIGINT
TimeoutSec=60s
Restart=on-failure
RestartSec=120s

[Install]
WantedBy=default.target
